# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:uwsgi]
command =/uwsgi/uwsgi --ini  /var/www/app/uwsgi.ini